# gentoo/overlay/sys-apps/systemd/files/update-etc-systemd-symlinks.path
[Unit]
Description=Update /etc/systemd/system symlinks for units moved to /usr

[Path]
PathChanged=/lib/systemd/system
